# pidacsd
#!/bin/bash

### BEGIN INIT INFO
# Provides:           PiDACS server
# Required-Start:     $local_fs $syslog $network
# Required-Stop:      $local_fs $syslog $network
# Default-Start:      2 3 4 5
# Default-Stop:       0 1 6
# Short-Description:  PiDACS server
# Description:        PiDACS server
### END INIT INFO

# Raspberry Pi Data Acquisition and Control System (PiDACS)

# PiDACS server daemon script (pidacsd).

cmd=$1
port_names=${2:-'gg0 gg1'}
usr=${3:-/usr/local}
bin=${usr}/bin

case $cmd in
    start)
        $bin/pidacs -L DATA "$port_names" 2>&1 & ;;
    stop)
        killall -q pidacs ;;
    restart)
        killall -q pidacs
        $bin/pidacs -l "$port_names" 2>&1 & ;;
    *)
        echo "pidacsd: usage: sudo /etc/init.d/pidacsd {start|stop|restart} ['port_names'] [usr dir]"
        exit 1 ;;
esac
